# supervisor.conf
[program:update_websites]
command=php /var/www/html/httpstatus/console.php controllers/internals/InternalApi update_websites
autostart=true
autorestart=true
